# crontab/delwc.sh
#!/bin/bash

# remove all wechat files
rm -rf $HOME/Library/Containers/com.tencent.xinWeChat/Data/Library/Application\ Support/com.tencent.xinWeChat/2.0b4.0.9
