# crontab/gitpull.sh
#!/bin/bash

# update repos on time
cd /path-to-your-repo; git pull
